# scripts/loop/com.juicyjusung.qgh-loop.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.juicyjusung.qgh-loop</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-lc</string>
        <string>exec /Users/user/projects/juicyjusung/juicy-qgh/scripts/loop/qgh-loop.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>10800</integer>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/Users/user/Library/Logs/qgh-loop.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/user/Library/Logs/qgh-loop.log</string>
</dict>
</plist>
